# sync/etc/init/life_ghost.conf
description "My life blog on Ghost"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [!12345]

# If the process quits unexpectadly trigger a respawn
respawn

setuid ghost
setgid ghost
env NODE_ENV=production
chdir /var/www/life_blog

exec /usr/local/bin/npm start --production

pre-stop exec /usr/local/bin/npm stop --production
